# src/battery_tester.ini
[program:battery-tester]
directory=/Users/mihaidobre/PycharmProjects/source/battery_tester/src/
command=/Users/mihaidobre/.virtualenvs/battery_tester/bin/python manage.py runserver 0.0.0.0:8000
autostart=true
autorestart=true
stderr_logfile=/Users/mihaidobre/PycharmProjects/source/battery_tester/src/logs/supervisor_err.log
stdout_logfile=/Users/mihaidobre/PycharmProjects/source/battery_tester/src/logs/supervisor.log

[program:celery-workers]
directory=/Users/mihaidobre/PycharmProjects/source/battery_tester/src/
command=/Users/mihaidobre/.virtualenvs/battery_tester/bin/celery worker --app=backend.celery:app -l info -c 5 --pool=eventlet -Q main_COM0,main_COM1,periodic_COM0,periodic_COM1,main_COM2,periodic_COM2,main_COM3,periodic_COM3
autostart=true
autorestart=true
stderr_logfile=/Users/mihaidobre/PycharmProjects/source/battery_tester/src/logs/workers_err.log
stdout_logfile=/Users/mihaidobre/PycharmProjects/source/battery_tester/src/logs/workers.log

[program:celery-beat]
directory=/Users/mihaidobre/PycharmProjects/source/battery_tester/src/
command=/Users/mihaidobre/.virtualenvs/battery_tester/bin/celery -A backend --app=backend.celery:app beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
autostart=true
autorestart=true
stderr_logfile=/Users/mihaidobre/PycharmProjects/source/battery_tester/src/logs/beat_err.log
stdout_logfile=/Users/mihaidobre/PycharmProjects/source/battery_tester/src/logs/beat.log
